# playbooks/roles/prometheus/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network-online.target

[Service]
ExecStart={{ node_exporter_dir }}/node_exporter-{{ node_exporter_version }}.{{ arch }}/node_exporter

Restart=always